# systemd/yeliztli-huey.service
[Unit]
Description=Yeliztli Huey Task Consumer
After=network.target yeliztli-api.service

[Service]
Type=simple
ExecStart=huey_consumer backend.tasks.huey_tasks.huey -w 1
WorkingDirectory=%h/Yeliztli
Restart=on-failure
RestartSec=5
Environment=PATH=%h/.local/bin:/usr/bin

[Install]
WantedBy=default.target
