# launchd/com.yeliztli.huey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yeliztli.huey</string>

    <key>ProgramArguments</key>
    <array>
        <string>huey_consumer</string>
        <string>backend.tasks.huey_tasks.huey</string>
        <string>-w</string>
        <string>1</string>
    </array>

    <!-- Replace __INSTALL_DIR__ with the actual Yeliztli install path -->
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>~/Library/Logs/yeliztli-huey.log</string>

    <key>StandardErrorPath</key>
    <string>~/Library/Logs/yeliztli-huey.log</string>
</dict>
</plist>
