# launchd/com.yeliztli.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yeliztli.api</string>

    <key>ProgramArguments</key>
    <array>
        <string>uvicorn</string>
        <string>backend.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8000</string>
    </array>

    <!-- Replace __INSTALL_DIR__ with the actual Yeliztli install path -->
    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>~/Library/Logs/yeliztli-api.log</string>

    <key>StandardErrorPath</key>
    <string>~/Library/Logs/yeliztli-api.log</string>
</dict>
</plist>
